# luci-app-campusnet/root/etc/uci-defaults/luci-campusnet
uci -q batch <<-EOF >/dev/null
	delete ucitrack.@campusnet[-1]
	add ucitrack campusnet
	set ucitrack.@campusnet[-1].init=campusnet
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
