# relevance/homebox/luci-app-transmission/root/etc/uci-defaults/luci-transmission

rm -rf /etc/seccomp/transmission-daemon.json

exit 0
