# luci-app-ssrust/root/etc/uci-defaults/40_luci-shadowsocks-rust
#!/bin/sh

mkdir -p /etc/shadowsocks-rust
/etc/init.d/rpcd reload
exit 0
