# luci-app-synology/root/etc/uci-defaults/luci-app-synology
#!/bin/sh

# if [ "$PKG_UPGRADE" != 1 ]; then
# 	uci batch <<-EOF
# 		set luci.themes.vzant=/luci-static/vzant
# 		set luci.main.mediaurlbase=/luci-static/vzant
# 		uci commit
# 	EOF
# fi

exit 0


