# luci-app-coredns/root/etc/hotplug.d/iface/99-coredns
#!/bin/sh
[ "$ACTION" = ifup ] && /etc/init.d/coredns restart
