# luci-theme-Light/root/etc/uci-defaults/luci-theme-Light
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Light=/luci-static/Light
	set luci.main.mediaurlbase=/luci-static/Light
	commit luci
EOF
exit 0
