# relevance/kwrt-packages/luci-app-sakurafrp/root/etc/uci-defaults/luci-app-sakurafrp
#!/bin/sh
[ ! -s "/etc/config/sakurafrp" ] && cp -f /usr/share/sakurafrp/0_default_config /etc/config/sakurafrp
