# relevance/openwrt-packages/luci-app-mmconfig/root/etc/uci-defaults/61-luci-app-mmconfig
#!/bin/sh -x

if [ "$PKG_UPGRADE" != "1" ]; then
	touch /etc/config/modemconfig
	uci set modemconfig.@modem[0]=modem
fi

