# relevance/my-default-settings/target/ath79/nand/base-files/etc/uci-defaults/LTE
if [[ "$(board_name)" == *xe300 ]]; then
uci -q add_list firewall.@zone[1].network='wwan0'
uci commit firewall

uci -q set network.wwan0=interface
uci -q set network.wwan0.proto='qmi'
uci -q set network.wwan0.device='/dev/cdc-wdm0'
uci -q set network.wwan0.apn='internet'
uci -q set network.wwan0.auth='none'
uci -q set network.wwan0.delay='10'
uci -q set network.wwan0.pdptype='IPV4V6'
uci commit network
fi
exit 0
