# luci-theme-opentomcat/applications/luci-app-docker/root/etc/uci-defaults/luci-docker
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@docker[-1]
	add ucitrack docker
	set ucitrack.@docker[-1].init=docker
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
