# luci-app-syscontrol/root/etc/uci-defaults/luci-syscontrol
#!/bin/sh

chmod +x /etc/init.d/ram_release /usr/bin/ram_release.sh

exit 0
