# relevance/my-default-settings/files/etc/hotplug.d/ieee80211/10-wifi-detect
#!/bin/sh

[ "${ACTION}" = "add" ] && {
	/sbin/wifi config
	[[ -z "$(uci -q get base_config.@status[0].wifi)" ]] && . /sbin/mywifi &
}
